# meteosat/etc/systemd/system/meteosat.timer
[Unit]
Description=Starts meteosat.service every hour

[Timer]
OnBootSec=10s
OnCalendar=*-*-* *:12,48:00

[Install]
WantedBy=timers.target
