# meteosat/DEBIAN/postinst
#!/bin/bash
systemctl start meteosat.timer
systemctl enable meteosat.timer
