# meteosat/DEBIAN/postinst
#!/bin/bash
systemctl enable meteosat.timer
